# deploy/systemd/openclaw-api.service
[Unit]
Description=OpenClaw API (uvicorn)
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/app
ExecStart=/home/openclaw/venv/bin/uvicorn apps.api.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/openclaw/.env

[Install]
WantedBy=multi-user.target
